# ansible-playbooks/ubuntu/20.04/server/suricata/roles/suricata/files/systemd/suricata-update.service
[Unit]
Description=Update suricata database

[Service]
Type=oneshot
ExecStart=/usr/local/bin/suricata-update

[Install]
WantedBy=default.target
