# systemd-user-setup/ARCHIVE/user/hdd-wakeup.service
[Unit]
Description=Wake up the external drive

[Service]
Type=oneshot
ExecStart=/bin/ls /media/XTORAGE-2TB
