# ansible_playbooks/ubuntu/18.04/server/initial-setup/roles/suricata/files/systemd/suricata-update.timer
[Unit]
Description=suricata-update timer

[Timer]
Persistent=true
OnCalendar=*-*-* 08:00:00
Unit=suricata-update.service

[Install]
WantedBy=timers.target
