# ansible_playbooks/ubuntu/18.04/initial-setup/roles/sysstat/files/systemd/sar-report-daily.service
[Unit]
Description=sar (sysstat) daily report

[Service]
Type=oneshot
ExecStart=/var/log/sysstat/sar-report-daily.sh

[Install]
WantedBy=default.target
